# cron/crontab
*/5 	*	*	*	*	/bin/ovr_wp_cron
*/15  * * * * /bin/rbb_cron
